# deploy/goban.service
[Unit]
Description=GoBan log watcher and IP banner
Documentation=https://github.com/izm1chael/goban
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
ExecStart=/usr/local/bin/goban-daemon --config /etc/goban/goban.yaml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s

# Capabilities — drop everything except what the banner needs.
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_DAC_OVERRIDE CAP_CHOWN

# systemd-managed directories with the right permissions.
# These create /run/goban, /var/lib/goban, /var/log/goban automatically.
RuntimeDirectory=goban
RuntimeDirectoryMode=0755
StateDirectory=goban
StateDirectoryMode=0700
LogsDirectory=goban
LogsDirectoryMode=0750

# Filesystem hardening
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
ReadWritePaths=/run/goban /var/lib/goban /var/log/goban
ReadOnlyPaths=/etc/goban

# Namespacing
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectControlGroups=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictNamespaces=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
